# doc/supervisor/horizon.ini
[program:horizon]
process_name=%(program_name)s
command=php /srv/http/mediadb/artisan horizon
autostart=true
autorestart=true
user=http
redirect_stderr=true
stdout_logfile=/path/for/horizon.log