# doc/supervisor/mediadb.ini
[program:horizon]
process_name=%(program_name)s
command=php /srv/http/mediadb/api/artisan horizon
autostart=true
autorestart=true
user=http
redirect_stderr=true
stdout_logfile=/srv/http/mediadb.log
stopwaitsecs=3600
