# doc/supervisor/websockets.ini
[program:websockets]
process_name=%(program_name)s
command=php /srv/http/mediadb/api/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=http
redirect_stderr=true
stdout_logfile=/srv/http/mediadb-websockets.log
stopwaitsecs=3600
