# doc/crontab
* * * * * cd /var/www/html/api && php artisan schedule:run >> /dev/null 2>&1
