# doc/etc/crontab
* * * * * cd /srv/http/mediadb/api && php artisan schedule:run >> /dev/null 2>&1
