# doc/supervisor/conf.d/horizon.ini
[program:horizon]
process_name=%(program_name)s
command=/usr/bin/php /var/www/html/api/artisan horizon
autostart=true
autorestart=true
user=http
redirect_stderr=true
stopwaitsecs=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
